# suspervisord.conf
[program:django]
command=uwsgi -T --enable-threads --process 4 --master --http :8000 --chdir /pythia-kore/src --module orax.wsgi
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=true